# scripts/com.chris.monday.clap-listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chris.monday.clap-listener</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
    <string>-u</string>
    <string>/Users/chris/CODE/MONDAY/scripts/double_clap_start_monday.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/chris/CODE/MONDAY</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>MONDAY_CLAP_AUDIO_DEVICE_INDEX</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/monday-clap-listener.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/monday-clap-listener.log</string>
</dict>
</plist>
